# telegram-bridge.service
[Unit]
Description=Paperclip Telegram Bridge Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/telegram-bridge
Environment=TELEGRAM_BOT_TOKEN=your-bot-token
Environment=PAPERCLIP_API_URL=http://127.0.0.1:3100
Environment=PAPERCLIP_COMPANY_ID=your-company-id
Environment=PAPERCLIP_API_KEY=your-api-key
Environment=TELEGRAM_CHAT_ID=your-chat-id
Environment=CEO_AGENT_ID=your-ceo-agent-id
Environment=HEARTBEAT_WEBHOOK_PUBLIC_ID=your-webhook-public-id
Environment=HEARTBEAT_WEBHOOK_SECRET=your-webhook-secret
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=node dist/index.js
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
